# php8.0-apache/fvolpato-php8.0-apache/crontab
* * * * * /usr/local/bin/php /var/www/html/persistMac.php >> /cron_scripts/mac_sync.log 2>&1

